# cloud/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:dagit]
environment=
  DAGSTER_HOME=/home/frsght/foresight/src/foresight,
  DAGSTER_DEPLOYMENT=production
command= /opt/conda/bin/conda run -n foresight /opt/conda/envs/foresight/bin/dagit -h 0.0.0.0 -w /home/frsght/foresight/src/foresight/workspace.yaml
autostart=true
autorestart=true
stderr_logfile = /var/log/supervisord/dagit-stderr.log
stdout_logfile = /var/log/supervisord/dagit-stdout.log

[program:dagster_daemon]
environment=
  DAGSTER_HOME=/home/frsght/foresight/src/foresight,
  DAGSTER_DEPLOYMENT=production
autostart=true
autorestart=true
command=/opt/conda/bin/conda run -n foresight /opt/conda/envs/foresight/bin/dagster-daemon run -w /home/frsght/foresight/src/foresight/workspace.yaml
stderr_logfile = /var/log/supervisord/dagster_daemon-stderr.log
stdout_logfile = /var/log/supervisord/dagster_daemon-stdout.log
